# src/integrationTest/docker/xenon-slurm/service/slurmctld/run
#!/bin/sh
exec /usr/sbin/slurmctld -D -c