# src/integrationTest/docker/xenon-slurm/service/slurmd-node-0/run
#!/bin/sh
exec /usr/sbin/slurmd -D -N node-0
